# crontab
0 7 * * 0-6 python /usr/src/app/onion-domain-harvester/main_harvester.py -d /usr/src/app/data >>/var/log/cron.log 2>&1
# Mandatory blank line