# crontab
0 8 * * * cd /usr/src/app/onion-domain-harvester && /usr/local/bin/python main_harvester.py 
# Mandatory blank line

